# docker/supervisord.conf
; supervisord pretends to be PID 1 inside the container. It runs the
; webserver (nginx + php-fpm) alongside the workers (queue + schedule)
; so a single Fly VM serves the whole app, including the cron-style
; pruners in routes/console.php.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue]
command=php artisan queue:work --tries=3 --sleep=2 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:schedule]
; schedule:work polls once a minute and dispatches due jobs. Matches
; the cron entries we'd otherwise put on the host: prune-expired-tokens
; (daily) and prune-activity-log (weekly).
command=php artisan schedule:work
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
